# deploy/systemd/clixer-analytics.service
[Unit]
Description=Clixer Analytics Service
Documentation=https://github.com/absadama/clixer_windows
After=network.target docker.service clixer-gateway.service
Requires=docker.service
StartLimitIntervalSec=60
StartLimitBurst=3

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/clixer/services/analytics-service
ExecStart=/usr/bin/npm run dev
Restart=on-failure
RestartSec=10

# Logging
StandardOutput=append:/opt/clixer/logs/analytics.log
StandardError=append:/opt/clixer/logs/analytics-error.log

# Environment
Environment=NODE_ENV=production
Environment=PORT=4005

# Process management
KillMode=mixed
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
