# deploy/systemd/clixer-etl-worker.service
[Unit]
Description=Clixer ETL Worker
Documentation=https://github.com/absadama/clixer_windows
After=network.target docker.service clixer-data.service
Requires=docker.service

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/clixer/services/etl-worker
ExecStart=/usr/bin/npm run dev
Restart=on-failure
RestartSec=10
StartLimitIntervalSec=60
StartLimitBurst=3

# Logging
StandardOutput=append:/opt/clixer/logs/etl-worker.log
StandardError=append:/opt/clixer/logs/etl-worker-error.log

# Environment
Environment=NODE_ENV=production

# Process management
KillMode=process
TimeoutStopSec=60

[Install]
WantedBy=multi-user.target
